# 00-bootstrap/dist/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>design.snds.workspace-doctor</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string><string>-lc</string>
    <string>WS="$(cat "$HOME/.claude/workspace-brain-path" 2>/dev/null)"; [ -f "$WS/AGENTS.md" ] || WS="$HOME/Projects/Workspace"; [ -f "$WS/AGENTS.md" ] || WS="$HOME/Projects/workspace"; "$WS/00-bootstrap/doctor/workspace-doctor.sh" --quiet</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>StartInterval</key><integer>14400</integer>
  <key>StandardErrorPath</key><string>/tmp/workspace-doctor.err</string>
</dict></plist>
